# trackpad-filter.service
[Unit]
Description=Trackpad dead zone filter
After=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/trackpad-filter --left 20 --right 20 --top 25
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
